# soundtrack-life-mvp/ops/systemd/qqmusic-api-web-restart.service
[Unit]
Description=Restart QQMusicApi Web Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart qqmusic-api-web.service
