# soundtrack-life-mvp/ops/systemd/soundtrack-life-proxy-restart.service
[Unit]
Description=Restart Soundtrack Life Backend Proxy

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart soundtrack-life-proxy.service
