# soundtrack-life-mvp/ops/systemd/soundtrack-life-proxy.service
[Unit]
Description=Soundtrack Life Backend Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/soundtrack-life
Environment=NODE_ENV=production
ExecStart=/usr/bin/node /opt/soundtrack-life/server/local-proxy.mjs
Restart=always
RestartSec=3
TimeoutStartSec=20
TimeoutStopSec=8
KillMode=control-group
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
